# deploy/contract_maker.service
# Systemd Service for Contract Maker
# Copy to: /etc/systemd/system/contract_maker.service
# Enable:  systemctl enable contract_maker
# Start:   systemctl start contract_maker
# Logs:    journalctl -u contract_maker -f

[Unit]
Description=Contract Maker - Leihvertrag Generator
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=contract_maker
Group=contract_maker
WorkingDirectory=/home/contract_maker/contract_maker

# Environment
Environment="FLASK_ENV=production"
Environment="KMS_MASTER_KEY_PATH=/etc/contract_maker/master.key"
EnvironmentFile=/home/contract_maker/contract_maker/.env

# Application
ExecStart=/home/contract_maker/contract_maker/venv/bin/gunicorn \
    --config gunicorn.conf.py \
    app:app
ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/contract_maker/contract_maker/users.db
ReadWritePaths=/home/contract_maker/contract_maker/static/signatures
ReadWritePaths=/home/contract_maker/contract_maker/signatures
ReadWritePaths=/var/log/contract_maker

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
